# vgpu.service
[Unit]
Description=lldp Service
After=network.target nvidia-vgpu-mgr.service multi-user.target

[Service]
Type=idle
RemainAfterExit=no

ExecStart=/etc/vcpu-vgpu/vcpu-vgpu.sh start
ExecStop=/etc/vcpu-vgpu/vcpu-vgpu.sh stop

[Install]
WantedBy=multi-user.target
